# polybar/.config/systemd/user/polybar-laptop.service
[Unit]
Description=Polybar

[Service]
ExecStart=/usr/bin/polybar laptop

[Install]
WantedBy=default.target
