# polybar/.config/systemd/user/polybar-side.service
[Unit]
Description=Polybar

[Service]
ExecStart=/usr/bin/polybar side

[Install]
WantedBy=default.target
