# services/.config/systemd/user/barrier.service
[Unit]
Description=Barrier
After=polybar.service
PartOf=polybar.service

[Service]
ExecStart=/usr/bin/barrier

[Install]
WantedBy=default.target
