# services/.config/systemd/user/pasystray.service
[Unit]
Description=pasystray
After=polybar.service
PartOf=polybar.service

[Service]
ExecStart=/usr/bin/pasystray

[Install]
WantedBy=default.target
