# picom/.config/systemd/user/picom.service
[Unit]
Description=picom
After=polybar.service
PartOf=polybar.service

[Service]
ExecStart=/usr/bin/picom

[Install]
WantedBy=default.target
